# infra/backup/bullpen-offsite.timer
[Unit]
Description=The Bullpen — offsite push timer (03:30 local, after the 03:00 snapshot)

[Timer]
# 30 minutes after the local snapshot timer; both are off-peak, outside any live-game
# window. Persistent=true catches up if the host was off when the timer should have fired.
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=10min
Unit=bullpen-offsite@alepic.service

[Install]
WantedBy=timers.target
